# deploy/relay/relay.service
[Unit]
Description=Network Next Relay

[Service]
ExecStart=/app/relay
EnvironmentFile=/app/relay.env
