# deploy/server-spawner.service
[Unit]
Description=Network Next Server Load Test

[Service]
ExecStart=/app/server-spawner.sh -n 30
Restart=always
RestartSec=3
